# usr/lib/systemd/system/pbhelper.service
[Unit]
Description=A Workaround for the power button problem.
After=phosh.service

[Service]
User=pbhelper
Group=pbhelper
ExecStart=/usr/bin/pbhelper

[Install]
WantedBy=graphical.target
